# etc/xdg/autostart/askilitetcxdgautostart.desktop
#!/usr/bin/env xdg-open
[Desktop Entry]
Exec=/usr/bin/askilitxdglink.sh
Type=Application
Name=askilitetcxdgautostart