# etc/xdg/autostart/askilitetcxdgautostart.desktop
#!/usr/bin/env xdg-open
[Desktop Entry]
Exec=sh /usr/bin/askilitxdglink.sh
Type=Application
Name[tr_TR]=askilitetcxdgautostart.desktop

